# backend/battleship.service
[Unit]
Description=Battleship WebSockets Server

[Service]
ExecStart=/usr/bin/lwsws -c /etc/battleship
StandardError=null

[Install]
WantedBy=multi-user.target

